# myproject.service
[Unit]
Description=Gunicorn instance to serve myproject
After=network.target

[Service]
User=asakovich
Group=www-data
WorkingDirectory=/home/asakovich/myproject
Environment="PATH=/home/asakovich/myproject/myprojectenv/bin"
ExecStart=/home/asakovich/myproject/myprojectenv/bin/gunicorn --workers 3 --bind unix:myproject.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target
